# deploy/templates/launchd/com.muzlife.backup-daily-db.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.muzlife.backup-daily-db</string>

  <key>WorkingDirectory</key>
  <string>__APP_ROOT__</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>./deploy/scripts/backup_daily_db.sh __APP_ROOT__</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__APP_ROOT__/runtime/logs/backup-daily-db.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__APP_ROOT__/runtime/logs/backup-daily-db.stderr.log</string>
</dict>
</plist>
