# runtime/plist/com.muzlife.qa-db-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.muzlife.qa-db-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/Caskroom/miniforge/base/bin/python3.12</string>
    <string>/Volumes/Data/Works/07.hahahoho/scripts/sync_prod_db_to_qa.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key><string>/Users/jingunpark/Library/Logs/hahahoho-qa/db-sync.out.log</string>
  <key>StandardErrorPath</key><string>/Users/jingunpark/Library/Logs/hahahoho-qa/db-sync.err.log</string>
</dict>
</plist>
