# runtime/plist/com.muzlife.library-qa.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>APP_ENV</key>
		<string>qa</string>
		<key>LIBRARY_DB_PATH</key>
		<string>/Volumes/Data/Works/07.hahahoho/data/library.db</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.muzlife.library-qa</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/Caskroom/miniforge/base/bin/python3.12</string>
		<string>-m</string>
		<string>uvicorn</string>
		<string>app.main:app</string>
		<string>--host</string>
		<string>0.0.0.0</string>
		<string>--port</string>
		<string>8100</string>
		<string>--app-dir</string>
		<string>/Volumes/Data/Works/07.hahahoho</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/jingunpark/Library/Logs/hahahoho-qa/launchd.qa.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/jingunpark/Library/Logs/hahahoho-qa/launchd.qa.stdout.log</string>
	<key>WorkingDirectory</key>
	<string>/Volumes/Data/Works/07.hahahoho</string>
</dict>
</plist>
